# crontab
*/5 * * * * cd /usr/app/news_crawl/news_crawl && /usr/local/bin/python main.py > /usr/app/news_crawl/crawler.log 2>&1
