# crontab
*/5 * * * * cd /usr/app/plw_rbtv77/plw_rbtv77 && /usr/local/bin/python main.py > /usr/app/plw_rbtv77/crawler.log 2>&1
